# .config/systemd/user/firefox-tmpfs-profile-sync@.timer
[Unit]
Description=Hourly synchronization of Firefox profile %I
StopWhenUnneeded=yes
RefuseManualStart=yes

[Timer]
OnCalendar=hourly
RandomizedDelaySec=5min

[Install]
WantedBy=default.target
