# .config/systemd/user/hdd-spinner@.timer
[Unit]
Description=Keep %I spinning

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=default.target
