# .config/systemd/user/mrss.service
[Unit]
Description=RSS to Mailbox

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mrss --config "%E/mrss/config"
