# .config/systemd/user/mrss.service
[Unit]
Description=RSS to Mailbox

[Service]
Type=oneshot
ExecStart=mrss --config "%E/mrss/config"
ExecStopPost=%E/mrss/post
OOMScoreAdjust=1000
