# .config/systemd/user/nvim-first-run.service
[Unit]
Description=Nvim First Run

[Service]
Type=oneshot
ExecStart=nvim --headless -c exit

[Install]
WantedBy=default.target
