# .config/systemd/user/sync.service
[Unit]
Description=Filesystem synchronization

[Service]
Type=simple
ExecStart=sync

[Install]
WantedBy=default.target
