# .config/systemd/user/sync.service
[Unit]
Description=Filesystem synchronization

[Service]
Type=oneshot
ExecStart=sync

[Install]
WantedBy=default.target
