# .config/systemd/user/sync.timer
[Unit]
Description=Filesystem Synchronization Timer

[Timer]
OnCalendar=00/8:00:00

[Install]
WantedBy=default.target
