# .config/systemd/user/sync.timer
[Unit]
Description=Filesystem Synchronization Timer

[Timer]
OnActiveSec=2h
OnUnitActiveSec=6h

[Install]
WantedBy=default.target
