# .config/systemd/user/sync.timer
[Unit]
Description=Hourly filesystem synchronization

[Timer]
OnCalendar=hourly

[Install]
WantedBy=default.target
